# BingWallpaper.cron
* */1 * * * sh ~/Pictures/BingWallpaper/bin/BingWallpaper.sh